# ansible-playbooks/roles/master/files/docker-registry.conf
[program:docker-registry]
directory=/tmp
user=root
command=/usr/bin/docker start --attach internal_registry
autostart=true
autorestart=true
